# test-interrupt/test-interrupt@.service
[Unit]
Description = test-interrupt@%i.service
PartOf = test-interrupt.target

[Service]
Type = oneshot
ExecStart = /usr/bin/sleep 30s

[Install]
WantedBy = test-interrupt.target
