# test-interrupt/test-instance.target
[Unit]
Description = test-instance.target 

[Install]
WantedBy = multi-user.target
